# system-root/etc/services.d/nginx/run
#!/usr/bin/execlineb -P

# Run nginx processis as wordpress:web
nginx -g "user wordpress web;"
